# stock_db.service
[Unit]
Description=Stock DB Service

[Service]
Type=simple
ExecStart=/usr/bin/stock_db.sh

[Install]
WantedBy=multi-user.target

